# services/cmc/Virtlyst/virtlyst
#!/sbin/openrc-run

depend() {
    need net virtlogd
    use virtlockd
    after libvirt
}

DAEMON=/usr/bin/cutelyst-wsgi2
DAEMON_OPTS="--ini /usr/Virtlyst/config.ini"

start() {
    ebegin "Starting virtlyst"
    start-stop-daemon -1 /var/log/virtlyst -2 /var/log/virtlyst --start --exec $DAEMON -b -- $DAEMON_OPTS
    eend $?
}

stop() {
     killall cutelyst-wsgi2
}
